# pi-dmx/engine/systemd/audio-dmx-engine.service
[Unit]
Description=Audio → DMX engine (Node)
After=network.target dmx-helper.service cpu-performance.service sound.target
Requires=dmx-helper.service

[Service]
Type=simple
WorkingDirectory=/opt/audio-dmx-engine
# --max-old-space-size caps V8 heap at 200 MB (Pi Zero 2 W has 512 MB total).
ExecStart=/usr/bin/node --enable-source-maps --max-old-space-size=200 dist/index.js
Environment=NODE_ENV=production
Environment=UV_THREADPOOL_SIZE=2
Environment=PORT=80
Restart=always
RestartSec=2

# Higher scheduling priority than default user processes, lower than dmx-helper's
# realtime SCHED_FIFO. Reduces audio-tick jitter under load.
Nice=-10
IOSchedulingClass=best-effort
IOSchedulingPriority=2

# Runs as root — this is a dedicated single-purpose appliance on an isolated
# network. Simplifies binding :80, ALSA access, and GPIO/serial permissions.
User=root

# Keep engine off the isolated realtime core (CPU3 reserved for dmx-helper).
CPUAffinity=0 1 2

[Install]
WantedBy=multi-user.target
