# pi-dmx/engine/systemd/pi-dmx-watchdog.service
[Unit]
Description=pi-dmx engine health watchdog (restarts engine if hung)
After=audio-dmx-engine.service

[Service]
Type=oneshot
ExecStart=/opt/audio-dmx-engine/watchdog.sh
