# pi-dmx/engine/systemd/pi-dmx-watchdog.timer
[Unit]
Description=Run the pi-dmx engine health watchdog every 30s

[Timer]
OnBootSec=60
OnUnitActiveSec=30
AccuracySec=5

[Install]
WantedBy=timers.target
